# ethernet_connect.service
[Unit]
Description=Ethernet Connection

[Service]
ExecStart=/home/aidanfoley/scripts/ethernet_connect.sh
